# distribution/image-builder.socket
[Unit]
Description=Image Builder socket

[Socket]
Service=image-builder.service
ListenStream=8086

[Install]
WantedBy=sockets.target
